# templates/start/upstart.conf.j2
# CI/CD managed
# Do NOT modify this file by hand!

description "{{item.upstart.description|default('The service')}}"

start on runlevel [2345]
stop on runlevel [016]

{% for key in item.upstart.env|default({}) %}
env {{key}}={{item.upstart.env[key]}}
{% endfor %}

chdir {{item.upstart.workingdirectory|default('/tmp')}}
setuid {{item.upstart.user|default('root')}}
exec {{item.upstart.execstart}}
{% for line in item.upstart.extralines|default([]) %}
{{ line|default('') }}
{% endfor %}
respawn
console log
